# proto/usr/lib/systemd/system/systemd-nspawn@.service.d/10-reboot-workaround.conf
#
#
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.
#

#
# Copyright 2020 Joyent, Inc.
#

#
# Remove --keep-unit from the command line so that the container can restart.
# https://github.com/systemd/systemd/issues/2809
#

[Service]
# First clear the existing ExecStart, else we end up with multiple commands.
ExecStart=
# Now define the new one
ExecStart=/usr/bin/systemd-nspawn --quiet --boot --link-journal=try-guest \
    --network-veth -U --settings=override --machine=%i
